# docker/scheduler/supervisor/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
