# docker/scheduler/supervisor/scheduler.conf
[program:scheduler]
process_name=%(program_name)s
directory=/var/www/html
command=php artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true

